# branches/release-0-22-fixes/99-mythadder.rules
SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/usr/bin/python /usr/bin/mythadder.py"
